# docker/etc/cron.d/redash
59 23 * * * root /tmp/bin/update_redash_default_date_range
